# torno-updater.service
[Unit]
Description=Actualizador OTA del torno QR
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/jesus/torno_qr/actualizador.py
User=root
